# bin/systemd/adsorber.timer
[Unit]
Description=Adsorber timer that runs the Adsorber service weekly.
Documentation=https://github.com/stablestud/adsorber

[Timer]
OnCalendar=weekly
Persistent=true
Unit=adsorber.service

[Install]
WantedBy=timers.target
